# Learn/Python/celery/03_Django/celery_learn/conf/supervisor_celery_beat.ini
[program:celery-beat]
command=python manage.py celery beat -l INFO
directory=/Users/leeyoung/PyCharmWorkstation/Learn/Python/celery/03Django/celery_learn
environment=PATH="/Users/leeyoung/Code/python3bottomvirtual/bin"
stdout_logfile=/Users/leeyoung/PyCharmWorkstation/Learn/Python/celery/03Django/celery_learn/logs/celery.beat.log
stderr_logfile=/Users/leeyoung/PyCharmWorkstation/Learn/Python/celery/03Django/celery_learn/logs/celery.beat.log
autostart=true
autorestart=true
startsecs=10
stopwatisecs=60
priority=997
