# learn/Python/celery/03_Django/celery_learn/conf/supervisor_celery_flower.ini
[program:celery-flower]
command=python manage.py celery flower
directory=/Users/leeyoung/PyCharmWorkstation/Learn/Python/celery/03Django/celery_learn
environment=PATH="/Users/leeyoung/Code/python3bottomvirtual/bin"
stdout_logfile=/Users/leeyoung/PyCharmWorkstation/Learn/Python/celery/03Django/celery_learn/logs/celery.flower.log
stderr_logfile=/Users/leeyoung/PyCharmWorkstation/Learn/Python/celery/03Django/celery_learn/logs/celery.flower.log
autostart=true
autorestart=true
startsecs=10
stopwatisecs=60
priority=996
